# buildroot-external/rootfs-overlay/etc/profile.d/resize_terminal.sh
#!/bin/sh

export SYSTEMD_COLORS=0
resize col 100 row 100  > /dev/null