# buildroot-external/rootfs-overlay/etc/profile.d/resize_terminal.sh
#!/bin/sh

resize col 100 row 100  > /dev/null